# roles/zfs/templates/etc/systemd/system/zfs-trim@.service.j2
[Unit]
Description=zpool trim for %i

[Service]
Nice=19
IOSchedulingClass=idle
KillSignal=SIGINT
ExecStart=zpool trim %i

[Install]
WantedBy=multi-user.target
