# roles/zfs/templates/etc/systemd/system/zfs-scrub@.service.j2
[Unit]
Description=zpool scrub for %i

[Service]
Nice=19
IOSchedulingClass=idle
KillSignal=SIGINT
ExecStart=zpool scrub %i

[Install]
WantedBy=multi-user.target
