# roles/zfs/templates/etc/systemd/system/zfs-scrub@.timer.j2
[Unit]
Description=zpool scrub for %i on calendar interval

[Timer]
OnCalendar={{ zfs_scrub_on_calendar }}
AccuracySec=1h
Persistent=true

[Install]
WantedBy=multi-user.target
