# roles/postfix_server/files/etc/systemd/system/postfix.service.d/override.conf
[Unit]
Wants=gssproxy.service
After=syslog.target network-online.target gssproxy.service

[Service]
Environment=GSS_USE_PROXY=yes
