# roles/dnf_automatic/templates/etc/systemd/system/dnf-automatic.service.d/override.conf.j2
{% if dnf_automatic_restart %}
[Service]
ExecStartPost={{ dnf_automatic_restart_script }}
{% endif %}
