# roles/asterisk/files/etc/systemd/system/asterisk.service.d/override.conf
[Unit]
After=nss-lookup.target network-online.target

[Service]
Restart=on-failure
RestartSec=4
