# roles/dnf_automatic/templates/etc/systemd/system/dnf-automatic.timer.d/override.conf.j2
[Timer]
OnCalendar={{ dnf_automatic_on_calendar }}
RandomizedDelaySec={{ dnf_automatic_random_delay }}
