# roles/apache/files/etc/systemd/system/httpd.service.d/override.conf
[Unit]
Wants=httpd-init.service gssproxy.service
After=network.target remote-fs.target nss-lookup.target httpd-init.service gssproxy.service

[Service]
Environment=GSS_USE_PROXY=yes
