# roles/dovecot/files/etc/systemd/system/dovecot.service.d/override.conf
[Unit]
Wants=gssproxy.service
After=local-fs.target network-online.target dovecot-init.service gssproxy.service

[Service]
Environment=GSS_USE_PROXY=yes
