# roles/postgresql_server/files/etc/systemd/system/postgresql.service.d/override.conf
[Unit]
Wants=gssproxy.service
After=network.target gssproxy.service

[Service]
Environment=GSS_USE_PROXY=yes
